# promet/setup/i386-linux/debian/postinst
#!/bin/bash
update-rc.d -f promet-process remove
update-rc.d promet-process defaults
